# supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/etc/supervisor/logs/supervisord.log
childlogdir=/etc/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock


[program:websocket-app]
process_name=websocket-app
command=php artisan serve --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1

redirect_stderr=true
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB

[program:websocket]
process_name=websocket
command=php artisan websocket:serve --port 31333
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true

stopwaitsecs=3600
stdout_logfile_maxbytes=5MB
